# tools/deploy/linux/systemd/pocket-option-activation-bot.service
[Unit]
Description=Pocket Option Telegram Activation Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__BOT_USER__
WorkingDirectory=__REPO_ROOT__
Environment=PYTHONPATH=src
EnvironmentFile=/etc/pocket-option-bot/bot.env
ExecStart=__REPO_ROOT__/.venv/bin/python tools/licensing/telegram_activation_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target